# celery.conf
#Celery task for grading-controller

description "Celery grading controller"
author "Vik Paruchuri <vik@edx.org>"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 3 30

env DJANGO_SETTINGS_MODULE=grading_controller.aws

chdir /opt/wwc/grading-controller
setuid www-data

exec /opt/edx/bin/python /opt/wwc/grading-controller/manage.py celeryd --loglevel=info --settings=grading_controller.aws --pythonpath=/opt/wwc/grading-controller -B --autoscale=4,1
